# deploy/vokaflow-backend.service
[Unit]
Description=VokaFlow Backend Service
Documentation=https://github.com/vokaflow/vokaflow
After=network.target redis.service postgresql.service
Wants=redis.service

[Service]
Type=simple
User=dw7
Group=vokaflow
WorkingDirectory=/opt/vokaflow
ExecStart=/opt/vokaflow/deploy/start-backend.sh
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vokaflow-backend

# Límites de recursos
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=4G
CPUQuota=200%

# Seguridad
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vokaflow /var/log/vokaflow /tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
